# openfang/rootfs/etc/services.d/openfang/run
#!/bin/sh
# Source all env vars written by cont-init.d (HOME, GITHUB_TOKEN, user env_vars, etc.)
for _f in /var/run/s6/container_environment/*; do
    [ -f "$_f" ] || continue
    _k=$(basename "$_f")
    _v=$(cat "$_f")
    export "${_k}=${_v}"
done
echo "[INFO] Starting openfang (listen: ${OPENFANG_LISTEN}, log: ${RUST_LOG})"
exec /usr/local/bin/openfang start
