# openfang/rootfs/etc/services.d/openfang/run
#!/bin/sh
echo "[INFO] Starting openfang (listen: ${OPENFANG_LISTEN}, log: ${RUST_LOG})"
exec /usr/local/bin/openfang start
