# openfang/rootfs/etc/s6-overlay/s6-rc.d/ttyd/finish
#!/command/execlineb -S0
foreground { s6-echo "ttyd exited" }
