# openfang/rootfs/etc/s6-overlay/s6-rc.d/openfang/run
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.info "Starting openfang (listen: ${OPENFANG_LISTEN}, log: ${RUST_LOG})"
exec /usr/local/bin/openfang start